# systemd/.config/systemd/user/gpg-agent.socket.d/override.conf
[Unit]
PartOf=smartcard.target

[Socket]
RemoveOnStop=true

[Install]
WantedBy=smartcard.target
