# systemd/.config/systemd/user/cvelist-sync.timer
[Unit]
Description=Sync cvelist repo hourly

[Timer]
OnBootSec=15m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
